# deploy/dashboard.service
[Unit]
Description=Personal Home Assistant Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Personal-Home-Assistant-Dashboard
EnvironmentFile=/home/pi/Personal-Home-Assistant-Dashboard/backend/.env
ExecStart=/home/pi/Personal-Home-Assistant-Dashboard/.venv/bin/uvicorn backend.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
